# rc.local
#!/bin/rc

ip link set wlan0 up &

wait
